# macos/com.personal-doctor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  Personal Doctor — macOS Launch Agent

  Auto-starts the health advisor server on login and keeps it running.
  Install with:
    cp macos/com.personal-doctor.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.personal-doctor.plist

  Uninstall with:
    launchctl unload ~/Library/LaunchAgents/com.personal-doctor.plist
    rm ~/Library/LaunchAgents/com.personal-doctor.plist

  View logs:
    tail -f ~/personal-doctor/logs/launchd-stdout.log
    tail -f ~/personal-doctor/logs/launchd-stderr.log
-->
<dict>
    <key>Label</key>
    <string>com.personal-doctor</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Updated by setup-mac.sh to point to the correct venv python -->
        <string>__VENV_PYTHON__</string>
        <string>-m</string>
        <string>app.server</string>
    </array>

    <key>WorkingDirectory</key>
    <!-- Updated by setup-mac.sh to the repo directory -->
    <string>__REPO_DIR__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Start on login, restart if it crashes -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Restart after a crash, wait 10 seconds -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Log stdout/stderr to files -->
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/launchd-stderr.log</string>
</dict>
</plist>
